# bin/vpnike
#!/bin/bash
### BEGIN INIT INFO
# Provides: ikec
# Required-Start: $network $syslog
# Required-Stop: $network $syslog
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Description: Start ikec
### END INIT INFO

SESSIONNAME=ikecsession
SITENAME=$PATH_FILE_IKE
VPNUSER=$USUARIO_VPN_IKE
VPNPASS=$SENHA_VPN_IKE

case "$1" in
start)
echo "Open ikec Connection..."
screen -S $SESSIONNAME -d -m ikec -r $SITENAME -u $VPNUSER -p $VPNPASS -a
;;

stop)
echo "Closing ikec Connection..."
 killall -9 screen
;;

*)
echo "error-->"
echo "Syntax: $0 {start|stop}"
exit 1
;;

esac
